# kai-services/streaming/units/kai-x11vnc.service
# リモート GUI（Tailscale 内のみ）。setup.sh が @TAILSCALE_IP@ を実 IP に置換して配置する。
[Unit]
Description=kai x11vnc on :0 (Tailscale only)
After=kai-desktop.service
Requires=kai-xorg.service

[Service]
ExecStart=/usr/bin/x11vnc -display :0 -listen @TAILSCALE_IP@ -no6 -rfbport 5900 -rfbauth %h/.vnc/passwd -forever -shared -noxdamage
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
